# dist/etc/init/solarsan-devserver.conf
# solarsan-devserver -- SolarSan Devserver
#
# ~trevorj 09/01/12

description	"SolarSan Devserver"

start on solarsan-start
stop on solarsan-stop

respawn
respawn limit 10 5

console log

env SOLARSAN_DIR=/opt/solarsanweb
env HOME=/root

chdir /opt/solarsanweb
#exec $SOLARSAN_DIR/manage runserver --werkzeug --dozer
#exec $SOLARSAN_DIR/manage runserver --werkzeug
exec $SOLARSAN_DIR/manage runserver
#command = /root/.virtualenvs/solarsanweb/bin/python $SOLARSAN_DIR/manage.py runserver_plus
#command = /root/.virtualenvs/solarsanweb/bin/python $SOLARSAN_DIR/manage.py runserver --werkzeug
